# systemd/formant-spot-adapter.service
# Example unit for packaged /opt installs.
# For a source checkout, scripts/deploy.sh writes the authoritative unit for the
# current machine and local config paths.
[Unit]
Description=Formant Spot Adapter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/formant-spot-adapter
ExecStart=/opt/formant-spot-adapter/scripts/run.sh
Restart=always
RestartSec=3
User=formant
Group=formant
Environment=ENV_FILE=/opt/formant-spot-adapter/config/formant-spot-adapter.env
Environment=CONFIG_PATH=/opt/formant-spot-adapter/config/formant-spot-adapter.json
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/opt/formant-spot-adapter/logs /opt/formant-spot-adapter/data
UMask=0077
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
